# init/annoy-diogo-and-nathan@.service
[Unit]
Description=Annoy diogo and nathan

[Service]
Type=oneshot
ExecStart=/home/jessie/annoy.sh
