# init/annoy-diogo-and-nathan@.timer
[Unit]
Description=Run script to annoy diogo and nathan

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
