# etc/systemd/user/dbus.socket
[Unit]
Description=D-Bus User Message Bus Socket

[Socket]
ListenStream=%t/bus

[Install]
WantedBy=sockets.target
Also=dbus.service
