# etc/update-motd.d/20-tail
#! /bin/sh
[ ! -f /etc/motd.tail ] && exit 0
cat /etc/motd.tail
