# etc/systemd/system/update-repos@.timer
[Unit]
Description=Run update-repos hourly and on boot

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
